# redhat-9/crontab
0 0 * * * /bin/bash /script/rsync-repo-redhat.sh