# redhat-9/crontab
0 0 * * * bash /script/rsync-repo-redhat.sh
