# redhat-8/crontab
0 0 * * * /bin/bash /script/rsync-repo-redhat.sh
0 1 * * * /bin/bash /script/rsync-repo-ubuntu-epel.sh